# vendors/InES/eva/rootfs/etc/inittab
#inet:unknown:/sbin/inetd -fe
#boa:unkown:/bin/boa
#udhcpd:unknown:/sbin/udhcpd -f
